# ops/oracle-bot.service
[Unit]
# Бот-пульт «Оракула» (раздел «Интерфейс»): пуш отчётов §8, кнопки решений §12,
# строка/алерты бюджета §15, алерты триггеров листа ожидания §17. Долгоживущий long-polling.
# Ревью 2026-07-04 Блок 4: работает от НЕПРИВИЛЕГИРОВАННОГО пользователя oracle (раньше root —
# любой дефект в обработке внешних данных давал root на хосте) + харднинг systemd.
Description=Oracle Telegram control panel (бот-пульт)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oracle
Group=oracle
WorkingDirectory=/home/oracle/oracle
# .env у проекта — строки `export KEY=val`; systemd EnvironmentFile их НЕ понимает (`export` ломает
# парсинг). Поэтому поднимаем окружение так же, как cron: bash + set -a + source .env, затем exec.
ExecStart=/bin/bash -c 'set -a && . /home/oracle/oracle/.env && set +a && exec /home/oracle/oracle/.venv/bin/python /home/oracle/oracle/ops/bot.py'
Restart=always
RestartSec=5
# Лог бота — рядом с cron.log проекта (systemd-журнал тоже пишется).
StandardOutput=append:/home/oracle/oracle/journal/bot.log
StandardError=append:/home/oracle/oracle/journal/bot.log
# Харднинг (ревью 04.07): без новых привилегий, свой /tmp, ядро/цгруппы только чтение,
# /usr /boot /etc — только чтение (проект в /home — пишется).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
